# gateway_proxy/files/gateway_proxy.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=54

start() {
	gateway_proxy &
	(/lib/functions/changeip.sh)
}

stop() {
	killall gateway_proxy
	sleep 1
}

